# config_files/supervisord-nginx.conf
[supervidord]
nodaemon=true

[program:nginx]
command=/opt/nginx/sbin/nginx -c /opt/nginx/conf/nginx.conf